# systemd/vega-undervolt-oc.service
[Unit]
Description=AMD Undervolt/Overclock Control

[Service]
ExecStart=/usr/local/bin/vega-undervolt-oc

[Install]
WantedBy=multi-user.target
